# scripts/infra/cloudflared.service
# ==============================================================================
# Cloudflare Tunnel Systemd Service
# ==============================================================================
# Exposes Strapi API to the internet via Cloudflare's edge network
#
# Prerequisites:
#   1. Install cloudflared: curl -L https://github.com/cloudflare/cloudflared/releases/latest/download/cloudflared-linux-arm64 -o /usr/local/bin/cloudflared
#   2. Create tunnel: cloudflared tunnel create godatify-api
#   3. Configure DNS: cloudflared tunnel route dns godatify-api api.godatify.com
#   4. Copy credentials: sudo cp ~/.cloudflared/<UUID>.json /etc/cloudflared/creds.json
#   5. Update /etc/cloudflared/config.yml with tunnel UUID
#
# Installation:
#   sudo cp cloudflared.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable cloudflared
#
# Usage:
#   sudo systemctl start cloudflared
#   sudo systemctl stop cloudflared
#   sudo systemctl status cloudflared
#   sudo journalctl -u cloudflared -f
# ==============================================================================

[Unit]
Description=Cloudflare Tunnel - godatify API
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/

# Start after Strapi is running (PM2 creates pm2-strapi service)
After=network-online.target pm2-strapi.service
Wants=network-online.target

# Optional: require Strapi to be running
# Wants=pm2-strapi.service

[Service]
# ------------------------------------------------------------------------------
# Process Configuration
# ------------------------------------------------------------------------------
Type=simple
User=cloudflared
Group=cloudflared

# ------------------------------------------------------------------------------
# Execution
# ------------------------------------------------------------------------------
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run

# Graceful reload (re-read config)
ExecReload=/bin/kill -HUP $MAINPID

# Graceful stop
TimeoutStopSec=10
KillMode=process

# ------------------------------------------------------------------------------
# Restart Policy
# ------------------------------------------------------------------------------
Restart=always
RestartSec=5

# Prevent restart loops
StartLimitIntervalSec=120
StartLimitBurst=5

# ------------------------------------------------------------------------------
# Resource Limits
# ------------------------------------------------------------------------------
# Cloudflared is lightweight, limit memory to prevent runaway
MemoryMax=256M
MemoryHigh=200M

# Low CPU priority (Strapi gets preference)
CPUWeight=50

# Open file limits
LimitNOFILE=8192

# ------------------------------------------------------------------------------
# Security Hardening
# ------------------------------------------------------------------------------
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Cloudflared only needs network access
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Only allow network access
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Restrict capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Read-only access to config
ReadOnlyPaths=/etc/cloudflared

# Deny write access everywhere else
ReadWritePaths=

# No privileges needed
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Lock down /proc and /sys
ProtectProc=invisible
ProcSubset=pid

# ------------------------------------------------------------------------------
# Logging
# ------------------------------------------------------------------------------
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared
SyslogLevel=info

[Install]
WantedBy=multi-user.target
